# backend/deploy/crawler-backend.service
[Unit]
Description=Crawler Backend Service
After=network.target

[Service]
User=crawler_user
WorkingDirectory=/opt/crawler/backend
Environment="PATH=/opt/crawler/venv/bin"
ExecStart=/opt/crawler/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target